# rootfs/etc/s6-overlay/s6-rc.d/server/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

# --- Directories ---
mkdir -p /data/.claudecode /tmp/claude-paste

# --- Ensure /tmp/claude-paste is auto-allowed in Claude Code settings ---
# Write to both the default (~/.claude) and CLAUDE_CONFIG_DIR locations to cover both cases
_merge_permission() {
  local f="$1"
  local dir
  dir="$(dirname "$f")"
  mkdir -p "$dir"
  [ -f "$f" ] || echo '{"permissions":{"allow":[],"deny":[]}}' > "$f"
  jq '.permissions.allow |= (. + ["Read(//tmp/claude-paste/**)"] | unique)' \
    "$f" > "${f}.tmp" && mv "${f}.tmp" "$f"
}
_merge_permission "/root/.claude/settings.json"
_merge_permission "/data/.claudecode/settings.json"

# Write image-paste convention so Claude Code understands [Image: path] references
cat > /data/.claudecode/CLAUDE.md << 'EOF'
# Image Paste Convention

When the user pastes an image into the terminal, it is saved to disk and the text
`[Image: /path/to/image.png]` is typed at the cursor.

When you see `[Image: /path/to/image.png]`, read the image directly from the specified path.
EOF

# --- Environment ---
export CLAUDE_CONFIG_DIR="/data/.claudecode"

# --- Symlinks in /root to match SSH terminal navigation ---
ln -sf /homeassistant /root/config
ln -sf /homeassistant /root/homeassistant
for d in share media ssl backup addons addon_configs; do
  [ -d "/$d" ] && ln -sf "/$d" "/root/$d" || true
done

cd /homeassistant || true

# --- Read add-on options ---
FONT_SIZE=$(bashio::config 'terminal_font_size' '14')
THEME=$(bashio::config 'terminal_theme' 'dark')
AUTO_UPDATE=$(bashio::config 'auto_update_claude' 'true')

export TERMINAL_FONT_SIZE="${FONT_SIZE}"
export TERMINAL_THEME="${THEME}"

# --- Register Home Assistant MCP (idempotent) ---
# Remove then re-add every start so SUPERVISOR_TOKEN (stable for container lifetime) is current.
claude mcp remove homeassistant 2>/dev/null || true
if claude mcp add-json homeassistant \
    "{\"command\":\"hass-mcp\",\"env\":{\"HA_TOKEN\":\"${SUPERVISOR_TOKEN}\",\"HA_URL\":\"http://supervisor/core\"}}" \
    -s user 2>&1 | bashio::log.info; then
  bashio::log.info "Home Assistant MCP registered"
else
  bashio::log.warning "Home Assistant MCP registration failed — HA tools will not be available"
fi

# --- Auto-update Claude Code (non-blocking, starts after server is up) ---
if bashio::var.true "${AUTO_UPDATE}"; then
  (
    sleep 5
    if timeout 120 claude update 2>&1 | bashio::log.info; then
      bashio::log.info "Claude Code updated successfully"
    else
      bashio::log.warning "Claude Code update failed or timed out (120s)"
    fi
  ) &
fi

# --- Launch server (s6 manages restarts on exit) ---
bashio::log.info "Starting Claude Code terminal on port 7681"
exec bun /app/server.ts
