# rootfs/etc/profile.d/claude-path.sh
export PATH="/root/.local/bin:/root/.bun/bin:${PATH}"
